# roles/hadoop/templates/hadoop.service.j2
[Unit]
Description=Hadoop Development service

[Service]
Type=oneshot
User={{ hadoop_owner }}
Group={{ hadoop_owner }}
Environment="JAVA_HOME={{ java_home_result.stdout }}"
Environment="{{ hadoop_base_path }}"
ExecStart={{ hadoop_base_path }}/sbin/start-all.sh
ExecStop={{ hadoop_base_path }}/sbin/stop-all.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target